# 0.9.1/debian/postinst
#!/bin/sh

set -e

echo "Enabling+Starting Service"

systemctl enable registrator
systemctl start registrator
