# 1.0.alpha.404/debian/prerm
#!/bin/sh

set -e

echo "Disabling Registrator Service"

systemctl stop registrator
systemctl disable registrator
